# .config/autostart/mount-cloud.desktop
[Desktop Entry]
Exec=/home/baudouin/bin/mount-cloud
Icon=dialog-scripts
Name=mount-cloud
Path=
Type=Application
X-KDE-AutostartScript=true
